# fulfillment/systemd/main-site-fulfillment.service
[Unit]
Description=MAIN-SITE Fulfillment Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/MAIN-SITE
Environment=NODE_ENV=production
EnvironmentFile=-/root/.openclaw/workspace/MAIN-SITE/fulfillment/.env
ExecStart=/usr/bin/npm run fulfillment:start
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
